# deploy/cdxtg.service
[Unit]
Description=cdxtg - Telegram control for Codex
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=CHANGE_ME
Group=CHANGE_ME
WorkingDirectory=/opt/cdxtg
Environment=NODE_ENV=production
Environment=CDXTG_ENV_FILE=/etc/cdxtg/telegram.env
ExecStart=/usr/bin/node --disable-warning=ExperimentalWarning /opt/cdxtg/dist/src/index.js
Restart=on-failure
RestartSec=5

# Codex applies the mode selected in Telegram. Do not add systemd sandboxing
# here: even apparently unrelated restrictions can imply NoNewPrivileges or a
# mount namespace and break danger-full-access or setgid tools such as crontab.

[Install]
WantedBy=multi-user.target
